# recipes-core/sensorhub/files/sensorhub-control.service
[Unit]
Description=Sensor Hub Control Service

[Service]
ExecStart=/opt/sensorhub/bin/sensorHubControl.lua 

[Install]
WantedBy=multi-user.target

