# recipes-core/silverline/files/silverline.target
[Unit]
Description=Silverline group of files
After=multi-user.target
Requires=multi-user.target

